# tests/layout/DEBIAN/postinst
#!/bin/sh
launchctl load /Library/LaunchDaemons/com.rpetrich.rocketbootstrap.tests.xpcserver.plist
exit 0
